# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:data-agent]
command=python -m main --agent data --host localhost --port 8001
directory=/app
autorestart=true
startretries=3
startsecs=10
priority=100
stdout_logfile=/var/log/data-agent.log
stderr_logfile=/var/log/data-agent-error.log

[program:planning-agent]
command=python -m main --agent planning --host localhost --port 8002
directory=/app
autorestart=true
startretries=3
startsecs=10
priority=101
stdout_logfile=/var/log/planning-agent.log
stderr_logfile=/var/log/planning-agent-error.log

[program:creative-agent]
command=python -m main --agent creative --host localhost --port 8003
directory=/app
autorestart=true
startretries=3
startsecs=10
priority=102
stdout_logfile=/var/log/creative-agent.log
stderr_logfile=/var/log/creative-agent-error.log

[program:host-agent]
command=python -m main --agent host --host localhost --port 8000
directory=/app
autorestart=true
startretries=3
startsecs=15
priority=200
stdout_logfile=/var/log/host-agent.log
stderr_logfile=/var/log/host-agent-error.log

[program:web-ui]
command=python run_ui.py --port 7860
directory=/app
autorestart=true
startretries=3
startsecs=30
priority=300
stdout_logfile=/var/log/web-ui.log
stderr_logfile=/var/log/web-ui-error.log
